# rootfs/etc/services.d/syslogd/run
#!/bin/sh
exec 2>&1
exec /sbin/syslogd -n -S -O -

